# deploy/systemd/stocktrade-background-update.timer
[Unit]
Description=Schedule StockTrade Background Latest Trade Day Update

[Timer]
# 默认按北京时间执行。交易日 16:30 触发首次更新，
# 若服务因“当日数据未就绪”以 TEMPFAIL(75) 退出，则交由 service 每 10 分钟重试直到成功。
OnCalendar=Mon..Fri *-*-* 16:30:00 Asia/Shanghai
Persistent=true
Unit=stocktrade-background-update.service

[Install]
WantedBy=timers.target
